# .config/systemd/user/xdg-desktop-portal.service.d/override.conf
[Service]
Environment="XDG_CURRENT_DESKTOP=gnome"
ExecStartPre=/usr/bin/dbus-update-activation-environment --systemd XDG_CURRENT_DESKTOP=gnome
